# cookbook/templates/default/solr_init.erb
#!/bin/sh
#
# Startup script for solr
#
# chkconfig: - 86 14
# processname: solr_glynx
# description: Rails application server
#
### BEGIN INIT INFO
# Provides: solr_glynx
# Required-Start: $local_fs $remote_fs $network
# Required-Stop: $local_fs $remote_fs $network
# Default-Start: 3
# Default-Stop: 0 1 2 4 5 6
# Short-Description: start and stop unicorn
### END INIT INFO
 
# Source function library.
. /etc/rc.d/init.d/functions
 
prog="solr"
RAILS_ROOT="<%= @install_path %>"
lockfile=${LOCKFILE-/var/lock/subsys/solr}
pidfile=$RAILS_ROOT/solr/pids/<%=@environment%>/sunspot-solr-<%=@environment%>.pid
SLEEPMSEC=100000
RETVAL=0
ENV=<%=@environment%>
 
start() {
    echo -n $"Starting $prog: "
    cd $RAILS_ROOT
    su - <%=@user%> -c "cd $RAILS_ROOT && chruby-exec 1.9.3 -- bundle exec rake sunspot:solr:start RAILS_ENV=${ENV}"
    RETVAL=$?
    echo
    [ $RETVAL = 0 ] || [ $RETVAL = 1 ] && touch ${lockfile}
    return 0 #$RETVAL
}

stop() {
    echo -n $"Stopping $prog: "
    su - <%=@user%> -c "cd $RAILS_ROOT && chruby-exec 1.9.3 -- bundle exec rake sunspot:solr:stop RAILS_ENV=${ENV}"
    RETVAL=$?
    echo
    [ $RETVAL = 0 ] && rm -f ${lockfile} ${pidfile}
}

rh_status() {
    status -p $pidfile java
} 

rh_status_q() {
    rh_status >/dev/null 2>&1
}

restart() {
    echo -n $"Restarting $prog: "
    stop
    start
    RETVAL=$?
    echo
}
 
# See how we were called.
case "$1" in
    start)
        rh_status_q && exit 0

        start
        ;;
    stop)       
        if ! rh_status_q; then
            rm -f $lockfile
            exit 0
        fi

        stop
        ;;
    status)
        rh_status
        RETVAL=$?
        if [ $RETVAL -eq 3 -a -f $lockfile ] ; then
            RETVAL=2
        fi
        ;;     
    restart)
        restart
        ;;
    *)
        echo $"Usage: $prog {start|stop|restart|condrestart|try-restart|force-reload|upgrade|reload|status|help|configtest}"
        RETVAL=2
esac
 
exit $RETVAL
